# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:fpm]
command=/etc/init.d/php7.0-fpm start

[program:mongo]
command=/usr/bin/mongod